# checkwifi.cron
MAILTO=""
*/3 * * * * root bash -x /home/pi/checkwifi.sh >>/tmp/cw.txt 2>&1 && rm -f /tmp/cw.txt
