# checkwifi.cron
MAILTO=""
*/5 * * * * root /home/pi/checkwifi.sh
